# cluster_img/patch-root/etc/profile.d/test_env.sh
#!/bin/bash

export PATH=$PATH:/opt/bin/:/opt/sbin/